# interception/systemd/system/udevmon.service
[Unit]
Description=udevmon service

[Service]
ExecStart=/usr/local/bin/udevmon

[Install]
WantedBy=default.target
